# debian/postrm
#!/bin/sh
set -e

# Maintainer script: postrm
# Disable service on purge (removal of symlinks already handled by debhelper but ensure).

if [ -d /run/systemd/system ] && command -v systemctl >/dev/null 2>&1; then
    case "$1" in
        purge)
            systemctl disable open-friendly-net-detection-client.service >/dev/null 2>&1 || true
            ;;
    esac
fi

#DEBHELPER#

exit 0
